# recipes_dependencies/sysd-mqtt/files/sysd-mqtt.service
[Unit]
Description=mqtt application launch Service

[Service]
ExecStart=/bin/bash /usr/bin/sysd-mqtt.sh

[Install]
WantedBy=multi-user.target
